# conf/supervisor.conf
[program:urlshorter]
command=/home/code/urlshorter/conf/gunicorn
user=code
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/code/urlshorter/logs/gunicorn-error.log